# dot_config/systemd/user/waybar-window-rewrite-autogen.service
[Unit]
Description=Auto-generate Waybar window rewrite icons
Wants=hyprland-session.target
Requires=waybar.service
After=waybar.service
After=hyprland-session.target
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 %h/.config/waybar/autogen_window_rewrite.py
Environment=OPENCODE_MODEL=
Environment=WAYBAR_REWRITE_AUTOGEN_LOG_LEVEL=INFO
Environment=PYTHONUNBUFFERED=1
