# dot_config/systemd/user/waybar-watcher.path
[Unit]
Description=Watch Waybar config for changes

[Path]
Unit=waybar-watcher.service
PathModified=%h/.config/waybar/config.jsonc
PathModified=%h/.config/waybar/shared_modules.jsonc
PathModified=%h/.config/waybar/style.css

[Install]
WantedBy=graphical-session.target
