# dot_config/systemd/user/waybar-watcher.path
[Unit]
Description=Watch Waybar config for changes

[Path]
Unit=waybar-watcher.service
PathModified=%h/.config/waybar

[Install]
WantedBy=default.target
